# deploy/systemd/qot-backup.timer
[Unit]
Description=Run the QOT backup every day

[Timer]
OnCalendar=*-*-* 02:15:00 UTC
Persistent=true
RandomizedDelaySec=10m

[Install]
WantedBy=timers.target
